# build/covid19tracker.info.conf
[program:covid19tracker.info]
command=/opt/covid19tracker.info/covid-19-tracker/deploy/start_gunicorn
directory=/opt/covid19tracker.info/covid-19-tracker/
autostart=true
autorestart=true